# medmo-deployment/medmo-api.service
[Unit]
Description=MedMO-4B API Server
After=network.target

[Service]
Type=simple
User=medmo
WorkingDirectory=/opt/medmo-api
Environment=PATH=/opt/medmo-api/venv/bin
Environment=PYTHONPATH=/opt/medmo-api
Environment=CUDA_VISIBLE_DEVICES=0
Environment=TRANSFORMERS_CACHE=/opt/medmo-api/cache
ExecStart=/opt/medmo-api/venv/bin/uvicorn medmo_api:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target